# basket-dashboard.service
[Unit]
Description=Basket Odds Dashboard
After=network.target

[Service]
Type=simple
User=pc
WorkingDirectory=/home/pc/Projects/basket-odd
EnvironmentFile=/home/pc/Projects/basket-odd/.env
Environment=DASHBOARD_PORT=5151
ExecStart=/home/pc/Projects/basket-odd/venv/bin/python run.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
